# dotfiles/home.config.systemd.user.cycle_wallpaper.service
[Unit]
Description=Cycle the wallpaper regularly

[Service]
ExecStart=/path/to/cycle_wallpaper.sh 1h
