# dotfiles/home.config.systemd.user.emacs.service
[Unit]
Description=Make emacs run quickly in the background.

[Service]
ExecStart=emacs -u username --fg-daemon
