# dotfiles/home.config.systemd.user.redshift.service
[Unit]
Description=Makes the screen red at night

[Service]
ExecStart=gtk-redshift -l 40:-103.5
